# scripts/adventure_bot.service
[Unit]
Description=MCADV - MeshCore Choose Your Own Adventure Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change 'pi' to your username if different
User=pi
WorkingDirectory=/home/pi/MCADV
# Basic usage – offline story trees only, no LLM, listen on channel_idx 1
ExecStart=/home/pi/MCADV/venv/bin/python3 /home/pi/MCADV/adventure_bot.py --baud 115200 --channel-idx 1 --announce
# Automatically restart if the bot crashes
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# -------------------------------------------------------------------------
# LLM OPTIONS – uncomment ONE of the blocks below to enable AI story generation
# -------------------------------------------------------------------------

# Option A: Ollama on a machine on your LAN (recommended for Pi Zero)
# Run Ollama on a separate PC/Pi with: ollama serve
# then pull a small model: ollama pull llama3.2:1b
# ExecStart=/home/pi/MCADV/venv/bin/python3 /home/pi/MCADV/adventure_bot.py \
#   --baud 115200 --channel-idx 1 --announce \
#   --ollama-url http://192.168.1.50:11434 --model llama3.2:1b

# Option B: Groq cloud API (free tier, needs internet)
# Sign up at https://console.groq.com and get an API key
# ExecStart=/home/pi/MCADV/venv/bin/python3 /home/pi/MCADV/adventure_bot.py \
#   --baud 115200 --channel-idx 1 --announce \
#   --groq-key gsk_YOUR_KEY_HERE

# Option C: OpenAI API (paid, needs internet)
# ExecStart=/home/pi/MCADV/venv/bin/python3 /home/pi/MCADV/adventure_bot.py \
#   --baud 115200 --channel-idx 1 --announce \
#   --openai-key sk_YOUR_KEY_HERE

# -------------------------------------------------------------------------
# DEBUG – add -d to see verbose output in journalctl
# -------------------------------------------------------------------------
# ExecStart=/home/pi/MCADV/venv/bin/python3 /home/pi/MCADV/adventure_bot.py \
#   --baud 115200 --channel-idx 1 -d

[Install]
WantedBy=multi-user.target
